# smartoys/smartoys.crontab
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/2 * * * * root  chmod -R +777 /var/www/app/logs
*/2 * * * * root  chmod -R +777 /var/www/app/cache
* 1 * * *   root  service cron reload
0,15,30,45 * * * *     www-data   php /var/www/html/app/console mautic:queue:process --env=prod -i page_hit > /var/log/cron.pipe 2>&1
0,15,30,45 * * * *     www-data   php /var/www/html/app/console mautic:queue:process --env=prod -i email_hit > /var/log/cron.pipe 2>&1

